# server/launchd/com.openrange.saxo-proxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.openrange.saxo-proxy</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>node</string>
    <string>/PATH/TO/REPO/server/index.js</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/PATH/TO/REPO/server</string>
  <key>StandardOutPath</key>
  <string>/PATH/TO/REPO/server/server.out.log</string>
  <key>StandardErrorPath</key>
  <string>/PATH/TO/REPO/server/server.err.log</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
</dict>
</plist>

<!--
Usage:
1. Copy this file to ~/Library/LaunchAgents/com.openrange.saxo-proxy.plist
2. Edit all /PATH/TO/REPO placeholders to the absolute path of your repository.
3. Load with: launchctl load ~/Library/LaunchAgents/com.openrange.saxo-proxy.plist
4. Unload with: launchctl unload ~/Library/LaunchAgents/com.openrange.saxo-proxy.plist
-->
